# containers/plausible-db.container
# containers/plausible-db.container
# PostgreSQL 16 for Plausible CE — rootless Podman quadlet

[Unit]
Description=Plausible PostgreSQL
After=network-online.target

[Container]
Image=docker.io/postgres:16-alpine
ContainerName=plausible-db
EnvironmentFile=%h/.config/containers/systemd/db.env
Network=plausible.network

Volume=plausible-db.volume:/var/lib/postgresql/data:Z

# Space-separated on one line — multiple AddCapability= lines also valid
# DropCapability=ALL must be uppercase; value is passed verbatim to --cap-drop
DropCapability=ALL
AddCapability=CHOWN FOWNER SETUID SETGID
NoNewPrivileges=false

Label=net.matrix.owner=FC13F74B@matrix.net
Label=net.matrix.project=plausible-quadlet-setup

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
