# containers/plausible.container
# containers/plausible.container
# Plausible CE analytics — rootless Podman quadlet
#
# Image pinned to tag + platform-specific digest (linux/amd64).
# Use platform manifest digest, NOT the OCI image index digest.
# Podman 5.x rejects index digests with "invalid checksum digest length".
# To get: skopeo inspect --raw docker://ghcr.io/plausible/community-edition:v3.2.1 | python3 -c "import sys,json; [print(m['digest']) for m in json.load(sys.stdin)['manifests'] if m.get('platform',{}).get('architecture')=='amd64']"
#
# Image pinned to tag + digest:
#   ghcr.io/plausible/community-edition:v3.2.1
#   sha256:7450d9df4bfce160541d65bdba6bd4bcdd9a6db07f13dde91060705fa242c650 (linux/amd64 platform manifest)
#   index: sha256:33e60bfb40f2df5da00f8753b76fad04f67dba3abe6d73eb516e440e3fb62985 (do not use — OCI index digest rejected by Podman 5.x)
#
# Dependency graph:
#   BindsTo= — if DB units stop, plausible stops and restarts when they return
#   Requires= — systemd will start the DBs when plausible starts
#   After= — enforces start ordering: DBs must be active before plausible starts

[Unit]
Description=Plausible CE analytics
Documentation=https://plausible.io/docs/self-hosting
After=network-online.target
After=plausible-db.service
After=plausible-events-db.service
Requires=plausible-db.service
Requires=plausible-events-db.service
BindsTo=plausible-db.service
BindsTo=plausible-events-db.service

[Container]
Image=ghcr.io/plausible/community-edition:v3.2.1@sha256:7450d9df4bfce160541d65bdba6bd4bcdd9a6db07f13dde91060705fa242c650
ContainerName=plausible
EnvironmentFile=%h/.config/containers/systemd/plausible.env
Network=plausible.network

Volume=plausible-data.volume:/var/lib/plausible:Z
PublishPort=127.0.0.1:8000:8000

DropCapability=ALL
NoNewPrivileges=true

Label=net.matrix.owner=FC13F74B@matrix.net
Label=net.matrix.project=plausible-quadlet-setup
Label=net.matrix.env=production

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
