# containers/plausible.container
# containers/plausible.container
# Plausible CE analytics — rootless Podman quadlet
# Installs to %h/.config/containers/systemd/

[Unit]
Description=Plausible CE analytics
Documentation=https://plausible.io/docs/self-hosting
After=network-online.target
After=plausible-db.service
After=plausible-events-db.service
Requires=plausible-db.service
Requires=plausible-events-db.service

[Container]
Image=ghcr.io/plausible/community-edition:v2.1.4
ContainerName=plausible
EnvironmentFile=%h/.config/containers/systemd/plausible.env
Network=plausible.network

Volume=plausible-data.volume:/var/lib/plausible:Z
PublishPort=127.0.0.1:8000:8000

DropCapability=ALL
NoNewPrivileges=true

Label=net.matrix.owner=FC13F74B@matrix.net
Label=net.matrix.project=plausible-quadlet-setup
Label=net.matrix.env=production

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
