# quadlets/plausible.container
# plausible.container — Podman quadlet for Plausible CE
# net.matrix label schema: owner FC13F74B@matrix.net
#
# Prerequisites:
#   zfs create storage/containers/plausible  (mountpoint /srv/plausible)
#   zfs create storage/users/plausible       (mountpoint /var/lib/plausible)
#   useradd -r -d /var/lib/plausible plausible
#   loginctl enable-linger plausible
#   Install quadlets to ~/.config/containers/systemd/

[Unit]
Description=Plausible CE analytics
Documentation=https://plausible.io/docs/self-hosting
After=network-online.target
After=plausible-db.service
After=plausible-events-db.service
Requires=plausible-db.service
Requires=plausible-events-db.service

[Container]
Image=ghcr.io/plausible/community-edition:v2.1.4
ContainerName=plausible
EnvironmentFile=/srv/plausible/config/plausible.env
Network=plausible.network

Volume=/srv/plausible/data:/var/lib/plausible:Z
PublishPort=127.0.0.1:8000:8000

Label=net.matrix.owner=FC13F74B@matrix.net
Label=net.matrix.project=plausible-quadlet-setup
Label=net.matrix.env=production

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
