# quadlets/plausible-events-db.container
# plausible-events-db.container — ClickHouse for Plausible CE

[Unit]
Description=Plausible ClickHouse
After=network-online.target

[Container]
Image=clickhouse/clickhouse-server:24.3.3.102-alpine
ContainerName=plausible-events-db
Network=plausible.network

Volume=/srv/plausible/clickhouse:/var/lib/clickhouse:Z
Volume=/srv/plausible/clickhouse-logs:/var/log/clickhouse-server:Z

Label=net.matrix.owner=FC13F74B@matrix.net
Label=net.matrix.project=plausible-quadlet-setup

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
